# monitor_hotplug/80-monitor_hotplug.rules
# camera

ACTION=="change", SUBSYSTEM=="drm", TAG+="monitor_hotplug"
TAG=="monitor_hotplug", RUN+="/usr/local/bin/monitor_hotplug"
